# deploy/dv-automission.service
# dv-automission.service — auto-record manual driving on AMI Manual confirm.
# Always-on watcher: subscribes /ami/mission and starts dv-manual.service (the
# sensors-only recorder) when Manual (index 0) is CONFIRMED on the AMI and the
# pipeline is not running; stops it when a real mission is confirmed. No
# autonomy, no commands to the car — read /ami/mission, drive a local unit.
#
# ENABLED at boot by the installer. Per uDV#189 the AMI emits an index only on
# CONFIRM (not dial position), so recording starts after one button press
# (select Manual + confirm), not on power-on alone. Opt out with
# `touch /etc/dv/norecord` (live) or `systemctl disable --now dv-automission`.
#
# Versioned in IFS08-DV-PIPELINE/deploy/ — install with
#   sudo ~/dv_ws/src/IFS08-DV-PIPELINE/deploy/install_dv_pipeline_service.sh

[Unit]
Description=ISC DV auto-record on AMI Manual selection (watches /ami/mission)
Documentation=https://github.com/isc-fs/IFS08-DV-PIPELINE
# Needs the uDV micro-ROS stream to see /ami/mission at all.
After=network.target microros-agent.service
Wants=microros-agent.service

[Service]
Environment=RMW_IMPLEMENTATION=rmw_cyclonedds_cpp
User=isc
ExecStart=/usr/local/bin/dv_automission.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
