# deploy/dv-automission.service
# dv-automission.service — auto-record manual driving on AMI Manual selection.
# Always-on watcher: subscribes /ami/mission and starts dv-manual.service (the
# sensors-only recorder) whenever Manual (index 0) is selected and the pipeline
# is not running; stops it when a real mission is dialled in. No autonomy, no
# commands to the car — read /ami/mission, drive a local systemd unit.
#
# ENABLED at boot by the installer, so "power on, dial at Manual" records with
# zero operator steps. Opt out with `touch /etc/dv/norecord` (honoured live) or
# `sudo systemctl disable --now dv-automission.service`.
#
# ⚠ Inert until isc-fs/IFS08-DV-uDV#189: it only fires if the uDV publishes
# /ami/mission with the ASMS off. Until then it watches and does nothing.
#
# Versioned in IFS08-DV-PIPELINE/deploy/ — install with
#   sudo ~/dv_ws/src/IFS08-DV-PIPELINE/deploy/install_dv_pipeline_service.sh

[Unit]
Description=ISC DV auto-record on AMI Manual selection (watches /ami/mission)
Documentation=https://github.com/isc-fs/IFS08-DV-PIPELINE
# Needs the uDV micro-ROS stream to see /ami/mission at all.
After=network.target microros-agent.service
Wants=microros-agent.service

[Service]
Environment=RMW_IMPLEMENTATION=rmw_cyclonedds_cpp
User=isc
ExecStart=/usr/local/bin/dv_automission.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
