# deploy/dv-automission.service
# dv-automission.service — auto-record manual driving whenever the car is on.
# Always-on watcher: triggers on the uDV sensor stream (/imu) being live — car
# powered + uDV up — and starts dv-manual.service (the sensors-only recorder)
# when the pipeline is not running; stops it when the pipeline runs (`dv race`
# takes over) or when the car powers off. No AMI, no autonomy, no commands to
# the car — read /imu, drive a local systemd unit. (Legacy "automission" name;
# trigger is now car-power, not mission — see the script header.)
#
# ENABLED at boot by the installer → power the car and it records, zero operator
# steps, no AMI. Opt out with `touch /etc/dv/norecord` (honoured live) or
# `sudo systemctl disable --now dv-automission.service`.
#
# No firmware dependency: /imu already publishes continuously, ASMS-off
# (confirmed uDV#189).
#
# Versioned in IFS08-DV-PIPELINE/deploy/ — install with
#   sudo ~/dv_ws/src/IFS08-DV-PIPELINE/deploy/install_dv_pipeline_service.sh

[Unit]
Description=ISC DV auto-record on AMI Manual selection (watches /ami/mission)
Documentation=https://github.com/isc-fs/IFS08-DV-PIPELINE
# Needs the uDV micro-ROS stream to see /ami/mission at all.
After=network.target microros-agent.service
Wants=microros-agent.service

[Service]
Environment=RMW_IMPLEMENTATION=rmw_cyclonedds_cpp
User=isc
ExecStart=/usr/local/bin/dv_automission.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
