# deploy/dv-record.service
# dv-record.service — race SENSOR rosbag recorder (replayable HIL source).
# Pulled in automatically by dv-pipeline.service (WantedBy below + enable),
# so every racing start (race boot, `dv race`, `dv restart`) records a bag of
# the raw sensor inputs (/lidar_points + /imu by default) — replayable back
# through the pipeline off-car; deliberately NOT the pipeline's own outputs.
# PartOf= propagates stop/restart from the pipeline: `dv stop` also stops
# the recorder, cleanly (SIGINT so rosbag2 finalizes metadata.yaml).
# The script waits for pipeline warmup (/dv/status) before recording, forces
# a best-effort QoS reader (the Hesai cloud is best-effort), and honours
# /etc/dv/norecord as an opt-out.
#
# Versioned in IFS08-DV-PIPELINE/deploy/ — install with
#   sudo ~/dv_ws/src/IFS08-DV-PIPELINE/deploy/install_dv_pipeline_service.sh

[Unit]
Description=ISC DV rosbag recorder (race telemetry)
Documentation=https://github.com/isc-fs/IFS08-DV-PIPELINE
After=dv-pipeline.service
PartOf=dv-pipeline.service

[Service]
User=isc
Environment=HOME=/home/isc
ExecStart=/usr/local/bin/dv_record.sh
# SIGINT (not SIGTERM) so ros2 bag record finalizes the bag metadata.
KillSignal=SIGINT
TimeoutStopSec=30
Restart=on-failure
RestartSec=15

[Install]
WantedBy=dv-pipeline.service
